# Terraform/photogallery.conf
[program:photogallery]
directory=/home/ubuntu/photogallery
command=/usr/bin/python3 /home/ubuntu/photogallery/main.py
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
